# conf/supervisor.conf
#写入下面这个文件
[program:blog]
command=/usr/local/bin/gunicorn wsgi -c /conf/gunicorn.conf.py
autostart=true
autorestart=true