# conf/supervisor.conf
#写入下面这个文件
[program:blog]
command=/usr/local/bin/gunicorn wsgi -c /var/www/blog/conf/gunicorn.conf.py
directory=/var/www/blog
autostart=true
autorestart=true